# instalador/novusio.service
# =============================================================================
# Arquivo de Serviço Systemd para Site Novusio
# Configurado para produção com segurança e monitoramento
# =============================================================================

[Unit]
Description=Site Novusio - Node.js Application
Documentation=https://github.com/seu-usuario/site-novusio-html
After=network.target network-online.target
Wants=network-online.target

[Service]
# Configurações básicas
Type=simple
User=USUARIO_LINUX
Group=USUARIO_LINUX
WorkingDirectory=CAMINHO_PROJETO
Environment=NODE_ENV=production
Environment=PORT=PORTA_SERVIDOR

# Comando de execução
ExecStart=/usr/bin/node server/server.js
ExecReload=/bin/kill -HUP $MAINPID

# Configurações de reinicialização
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Configurações de timeout
TimeoutStartSec=60s
TimeoutStopSec=30s

# Configurações de logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=novusio

# Configurações de segurança
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=CAMINHO_PROJETO
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

# Configurações de recursos
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=80%

# Configurações de ambiente de trabalho
WorkingDirectory=CAMINHO_PROJETO
Environment=HOME=/home/USUARIO_LINUX
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
