# instalador/systemd.service
# =============================================================================
# SERVIÇO SYSTEMD PARA NOVUSIO
# =============================================================================
# Este arquivo configura o serviço systemd para inicialização automática
# do PM2 e da aplicação Novusio
# =============================================================================

[Unit]
Description=Novusio Application (PM2)
Documentation=https://pm2.keymetrics.io/
After=network.target

[Service]
Type=forking
User=novusio
Group=novusio
WorkingDirectory=/opt/novusio
ExecStart=/usr/local/bin/pm2 start /opt/novusio/ecosystem.config.js --env production
ExecReload=/usr/local/bin/pm2 reload all
ExecStop=/usr/local/bin/pm2 kill
Restart=always
RestartSec=10

# Configurações de segurança
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/novusio /var/log/novusio /opt/backups/novusio

# Configurações de recursos
LimitNOFILE=65536
LimitNPROC=4096

# Configurações de ambiente
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
